# hk-callback/feishu-callback.service
[Unit]
Description=Feishu Callback (OKX Pump Monitor) - HK self-hosted
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=feishu
Group=feishu
WorkingDirectory=/opt/feishu-callback
EnvironmentFile=/opt/feishu-callback/.env
# 监听 127.0.0.1:8000，外部由 Caddy 反代过来
ExecStart=/opt/feishu-callback/venv/bin/uvicorn app:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1 \
    --log-level info \
    --access-log
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
# 安全沙箱
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/feishu-callback
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
